# system_files/etc/profile.d/neofetch.sh
alias neofetch='/usr/bin/neowofetch --config /usr/share/ublue-os/bazzite/neofetch.conf'
alias neowofetch='/usr/bin/neowofetch --config /usr/share/ublue-os/bazzite/neofetch.conf'
